# ipadclaw.service
[Unit]
Description=iPadClaw HTTP Bridge for OpenClaw
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/csabaszabo/ipadclaw
ExecStart=/home/linuxbrew/.linuxbrew/bin/node /home/csabaszabo/ipadclaw/server.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=XDG_RUNTIME_DIR=/run/user/1000

[Install]
WantedBy=default.target
